# server/ansible/roles/goaccess/templates/visitorbot.service
[Unit]
Description=Outputs visitor information into the Zulip team chat

[Service]
Type=oneshot
ExecStart=/usr/bin/visitorbot
