# server/ansible/roles/goaccess/templates/visitorbot.timer
[Unit]
Description=Run visitorbot.service when the day changes

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
